# setup/profile.d/pywal.sh
#!/bin/sh

wal -i ~/wallpapers -steq" >> /etc/profile